# mkosi.profiles/embedded-firmware/mkosi.extra/usr/lib/systemd/system/sdm845-fwload.service.d/embedded-firmware.conf
[Unit]
ConditionPathExists=!/usr/lib/firmware/.embedded-firmware
