# mkosi.profiles/sdm845-fwload/mkosi.extra/usr/lib/systemd/system/iio-sensor-proxy.service.d/wait-for-hexagonrpcd.conf
[Unit]
After=hexagonrpcd-sdsp.service
Requires=hexagonrpcd-sdsp.service
