# mkosi.profiles/droid-juicer/mkosi.extra/usr/lib/systemd/system/phrog.service.d/wait-for-droid-juicer.conf
[Unit]
After=droid-juicer.service
Requires=droid-juicer.service
